# C CODE/com.visualstudio.code.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
<key>Label</key>
<string>com.visualstudio.code.tunnel</string>
<key>LimitLoadToSessionType</key>
<string>Aqua</string>
<key>ProgramArguments</key>
<array>
<string>/Applications/Visual Studio Code.app/Contents/Resources/app/bin/code-tunnel</string>
<string>--verbose</string><string>--cli-data-dir</string><string>/Users/harshilpatel/.vscode/cli</string><string>tunnel</string><string>service</string><string>internal-run</string>
</array>
<key>KeepAlive</key>
<true/>
<key>StandardErrorPath</key>
<string>/Users/harshilpatel/.vscode/cli/tunnel-service.log</string>
<key>StandardOutPath</key>
<string>/Users/harshilpatel/.vscode/cli/tunnel-service.log</string>
</dict>
</plist>